# scribe/scribe.cron

* 5 * * * $SCRIBEUSER /usr/local/bin/scribebak
